# aml_linux_fd628/system.d/fd628.service
[Unit]
Description=FD628 service

[Service]
Type=simple
ExecStart=/usr/bin/fd628 &

[Install]
WantedBy=basic.target
